# supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=/usr/local/bin/docker-entrypoint.sh postgres

[program:app]
command=python3 /opt/api/app.py
